# examples/backend.container
[Container]
Image=backend
Pod=my-stack.pod
Volume=db-data:/var/lib/db
