# examples/frontend.container
[Container]
Image=frontend
Pod=my-stack.pod
Exec=/app/start
